# daten/persistenz/lms-persist-save.service
[Unit]
Description=Linux Mint Snapshot — RAM-Persistenz beim Herunterfahren speichern
After=basic.target
Before=umount.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/true
ExecStop=/usr/local/sbin/lms-persist-save
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
