# daten/scripts/rikus-mintshot-sshkeys.service
[Unit]
Description=Rikus Mintshot — SSH-Rechnerschluessel im Klon sicherstellen
Documentation=https://github.com/Zahnschmerz/rikus-mintshot
# Der Schnappschuss laesst die SSH-Rechnerschluessel bewusst weg (sonst haetten alle
# Klone dieselben). Beim LIVE-Start erzeugt live-config sie neu — nach einer FESTEN
# INSTALLATION aber NICHT. Dann startet sshd gar nicht ("no hostkeys available"),
# und zwar still: erst beim Zugriff von aussen faellt es auf.
# Dieser Dienst erzeugt fehlende Schluessel und startet sshd nach. Auf Systemen, die
# ihre Schluessel haben, tut er nichts.
# Reihenfolge WICHTIG: NICHT vor network-pre.target und NICHT mit DefaultDependencies=no
# starten! Im Live-Betrieb erzeugt Mints live-config die SSH-Schluessel selbst
# (/lib/live/config/1160-openssh-server). Lief dieser Dienst FRUEHER als live-config,
# machte er die Schluessel selbst und rief synchron "systemctl restart ssh" -> das
# wartete aufs Netz, das seinerseits auf uns wartete = Deadlock, der Live-Boot blieb
# haengen (kein Desktop). Normal eingereiht laeuft er NACH live-config: die Schluessel
# sind schon da -> er tut nichts. Fest installierte Systeme (kein live-config) versorgt
# er wie bisher, nur eingereiht vor ssh.service.
After=local-fs.target live-config.service
Before=ssh.service sshd.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/sbin/rikus-mintshot-sshkeys

[Install]
WantedBy=multi-user.target
